# upgrade-watcher.service
[Unit]
Description=Cosmos off-chain upgrade watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=/etc/upgrade-watcher.env
Environment=ENV_FILE=/etc/upgrade-watcher.env
RuntimeDirectory=upgrade-watcher
RuntimeDirectoryMode=0700
UMask=0077
ExecStart=/usr/local/bin/upgrade-watcher.sh
Restart=no

[Install]
WantedBy=multi-user.target
